# caltoc_startup.conf
description "Challenge Accepted Server"
author      "Zack Harris and Dante Camarena"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/zharris/"
	
    exec sudo /usr/local/bin/node /home/zharris/website/app.js >> /var/log/node.log 2>&1
end script

#post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   exec echo "node.js has started!" /home/zharris/website/notify_log
#end script
